# debian/cron.daily
#!/bin/sh

MAX_TMPFILE_LIFETIME=5
# MAX_TMPFILE_LIFETIME can be overridden in this file
if [ -r /etc/default/multiflexi-core ]; then
	. /etc/default/multiflexi-core
fi

if [ -d /var/lib/multiflexi/temp ]; then
	find /var/lib/multiflexi/temp -type f -mtime +$MAX_TMPFILE_LIFETIME -delete
fi
