# debian/multiflexi.cron.weekly
#!/bin/sh
sudo -u www-data multiflexi-executor w
