# debian/multiflexi.cron.d
#Run AbraFlexi scripts once per year
0 0 1 1 *	www-data multiflexi-scheduler y
#Schedule ad-hoc actions
* * * * *	www-data multiflexi-scheduler i
