# debian/multi-flexibee-setup.cron.d
#Run FlexiBee scripts once per year
0 0 1 1 *	www-data multi-flexibee-executor y

